# overlay/etc/supervisor/conf.d/py3o-fusion.conf
[program:py3o-fusion]
command=/usr/local/bin/start-py3o-fusion --renderserver=localhost
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
